# scripts/weibo-monitor.service
[Unit]
Description=岁己SUI 微博监控
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%u
WorkingDirectory=/opt/weibo-monitor
ExecStart=/opt/weibo-monitor/venv/bin/python3 /opt/weibo-monitor/app.py
EnvironmentFile=/opt/weibo-monitor/.env
Environment=PYTHONUNBUFFERED=1

# 自动重启（异常退出5秒后）
Restart=always
RestartSec=5

# 资源限制（e2-micro 友好）
CPUWeight=50
Nice=10
MemoryMax=128M

# 日志输出到 journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weibo-monitor

[Install]
WantedBy=multi-user.target
